# droneaware-ble.service
[Unit]
Description=DroneAware BLE Remote ID detector (offline)
After=bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pi/drone-aware-zero/ble_feeder.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
